# snapmaker/rootfs/etc/s6-overlay/s6-rc.d/smMonitor/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Snapmaker Monitor
# Runs Snapmaker Monitor script
# ==============================================================================

bashio::log.info "Starting Snapmaker Monitor..."

exec python3 /app/smStatus.py
