# snapmaker-monitor/rootfs/etc/s6-overlay/s6-rc.d/smMonitor/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Snapmaker Monitor
# Runs Snapmaker Monitor script
# ==============================================================================

bashio::log.info "Starting Snapmaker Monitor..."

if bashio::config.has_value 'log_level'; then
    log_level=$(bashio::config 'log_level')
    export LOG_LEVEL=${log_level}
fi

if bashio::config.has_value 'ha_webhook_id'; then
    ha_webhook_id=$(bashio::config 'ha_webhook_id')
    export HA_WEBHOOK_ID=${ha_webhook_id}
fi

if bashio::config.has_value 'sm_ip'; then
    sm_ip=$(bashio::config 'sm_ip')
    export SM_IP=${sm_ip}
fi

if bashio::config.has_value 'sm_port'; then
    sm_port=$(bashio::config 'sm_port')
    export SM_PORT=${sm_port}
fi

if bashio::config.has_value 'mqtt_broker'; then
    mqtt_broker=$(bashio::config 'mqtt_broker')
    export MQTT_BROKER=${mqtt_broker}
fi

if bashio::config.has_value 'mqtt_port'; then
    mqtt_port=$(bashio::config 'mqtt_port')
    export MQTT_PORT=${mqtt_port}
fi

if bashio::config.has_value 'mqtt_topic'; then
    mqtt_topic=$(bashio::config 'mqtt_topic')
    export MQTT_TOPIC=${mqtt_topic}
fi

if bashio::config.has_value 'mqtt_user'; then
    mqtt_user=$(bashio::config 'mqtt_user')
    export MQTT_USER=${mqtt_user}
fi

if bashio::config.has_value 'mqtt_password'; then
    mqtt_password=$(bashio::config 'mqtt_password')
    export MQTT_PASSWORD=${mqtt_password}
fi

exec python3 /app/smStatus.py
